# proxy.service.in
[Unit]
Description=Proxy Server Core Daemon
After=network.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/proxy_server --config /etc/proxy_server/%i.conf

LimitNOFILE=900000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
